# dockerfiles/lab-soc-analyst/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
autostart=true
autorestart=true
stdout_logfile=/var/log/elasticsearch-stdout.log
stderr_logfile=/var/log/elasticsearch-stderr.log
environment=ES_JAVA_OPTS="-Xms512m -Xmx512m"

[program:kibana]
command=/usr/share/kibana/bin/kibana
user=kibana
autostart=true
autorestart=true
startsecs=30
stdout_logfile=/var/log/kibana-stdout.log
stderr_logfile=/var/log/kibana-stderr.log

[program:wazuh-manager]
command=/var/ossec/bin/wazuh-control start
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/wazuh-stdout.log
stderr_logfile=/var/log/wazuh-stderr.log

[program:suricata]
command=/usr/bin/suricata -c /etc/suricata/suricata.yaml -i eth0 --init-errors-fatal
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/suricata-stdout.log
stderr_logfile=/var/log/suricata-stderr.log

[program:filebeat]
command=/usr/share/filebeat/bin/filebeat -e -c /etc/filebeat/filebeat.yml
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/var/log/filebeat-stdout.log
stderr_logfile=/var/log/filebeat-stderr.log

[program:ttyd]
command=/usr/local/bin/ttyd -p 7681 login -f lab
autostart=true
autorestart=true
stdout_logfile=/var/log/ttyd.log
stderr_logfile=/var/log/ttyd-err.log
