# dockerfiles/lab-bigdata-workspace/supervisord.conf
; =============================================================================
; Supervisord config for lab-bigdata-workspace
;
; Services are split into "always on" and "opt-in" groups. Opt-in services
; are started only if their corresponding env var is set, via the start.sh
; wrapper that writes supervisor programs dynamically before supervisord boots.
;
; ttyd (browser terminal) is always on. MySQL is always on. SSH, Kafka,
; Jupyter, Cassandra, and Spark master/worker are opt-in via env vars.
; =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

; -------------------------- Always-on services -----------------------------

[program:ttyd]
command=/usr/local/bin/ttyd -p 7681 -W -t fontSize=14 -t 'theme={"background":"#0a0a0a"}' login -f lab
autostart=true
autorestart=true
stdout_logfile=/var/log/ttyd.stdout.log
stderr_logfile=/var/log/ttyd.stderr.log
priority=10

[program:mysql]
command=/usr/sbin/mysqld --user=mysql --datadir=/var/lib/mysql --bind-address=127.0.0.1
autostart=true
autorestart=true
stdout_logfile=/var/log/mysql.stdout.log
stderr_logfile=/var/log/mysql.stderr.log
priority=20

; -------------------------- Opt-in services --------------------------------
; These are enabled/disabled by the start.sh wrapper writing their autostart
; value based on env vars (ENABLE_SSH, ENABLE_KAFKA, ENABLE_CASSANDRA,
; ENABLE_SPARK, ENABLE_JUPYTER). Defaults: SSH off, Kafka on, Cassandra off,
; Spark on, Jupyter off.

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=false
autorestart=true
stdout_logfile=/var/log/sshd.stdout.log
stderr_logfile=/var/log/sshd.stderr.log
priority=30

[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/kraft/server.properties
autostart=false
autorestart=true
environment=KAFKA_HEAP_OPTS="-Xmx512m -Xms256m"
stdout_logfile=/var/log/kafka.stdout.log
stderr_logfile=/var/log/kafka.stderr.log
priority=40

[program:spark-master]
command=/opt/spark/sbin/start-master.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/var/log/spark-master.stdout.log
stderr_logfile=/var/log/spark-master.stderr.log
priority=50

[program:spark-worker]
command=bash -c "sleep 5 && /opt/spark/sbin/start-worker.sh spark://localhost:7077"
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/var/log/spark-worker.stdout.log
stderr_logfile=/var/log/spark-worker.stderr.log
priority=51

[program:cassandra]
command=/opt/cassandra/bin/cassandra -f -R
user=cassandra
autostart=false
autorestart=true
environment=CASSANDRA_HEAP_NEWSIZE="200M",MAX_HEAP_SIZE="1G"
stdout_logfile=/var/log/cassandra.stdout.log
stderr_logfile=/var/log/cassandra.stderr.log
priority=60

[program:jupyter]
command=su - lab -c "jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --NotebookApp.token='' --NotebookApp.password=''"
autostart=false
autorestart=true
stdout_logfile=/var/log/jupyter.stdout.log
stderr_logfile=/var/log/jupyter.stderr.log
priority=70
